# deploy/pestcrm-site.service
[Unit]
Description=Drudge marketing website (static) on port 8080
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pest-crm/marketing-site
ExecStart=/usr/bin/python3 /root/pest-crm/marketing-site/server.py 8080
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
